# deploy/polyx-worker.service
[Unit]
Description=PolyX Worker (Copy Engine)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/polyx_bot_new
ExecStart=/usr/bin/python3 -m polyx_bot.worker
Restart=always
RestartSec=5
EnvironmentFile=/home/ubuntu/polyx_bot_new/.env

[Install]
WantedBy=multi-user.target
